# systemfiles/systemd/system/beamerdetect.service
[Unit]
Description=lockbutton
After=getty.target

[Service]
ExecStart=/root/pidor/scripts/beamerdetect.sh
Restart=always

[Install]
WantedBy=multi-user.target
