# systemfiles/systemd/system/caststatus.service
[Unit]
Description=lockbutton
After=getty.target

[Service]
ExecStart=/root/pidor/scripts/caststatus.py
Restart=always

[Install]
WantedBy=multi-user.target
